# rc.local
# example rc.local to run at startup
/root/operator.py &
exit 0
